# systemd/nextmeeting.service
[Unit]
Description=NextMeeting Server
After=graphical-session.target

[Service]
Type=simple
ExecStart=nextmeeting server
Restart=on-failure
RestartSec=5

[Install]
WantedBy=graphical-session.target
